# software/systemd/phenocam-capture.timer
[Unit]
Description=PhenoCam capture timer

[Timer]
OnCalendar=*:11,41
AccuracySec=10s
Persistent=false

[Install]
WantedBy=timers.target
